# crontab
PATH=/usr/local/bin
*/20 * * * * ( cd /icms_crawler/ && python3 main.py >> /var/log/icms_crawler.log 2>&1)
#
